# opends/src/svr4/OpenDS/opends-manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License, Version 1.0 only
 (the "License").  You may not use this file except in compliance
 with the License.

 You can obtain a copy of the license at
 trunk/opends/resource/legal-notices/OpenDS.LICENSE
 or https://OpenDS.dev.java.net/OpenDS.LICENSE.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at
 trunk/opends/resource/legal-notices/OpenDS.LICENSE.  If applicable,
 add the following below this CDDL HEADER, with the fields enclosed
 by brackets "[]" replaced with your own identifying information:
      Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END


      Copyright 2008 Sun Microsystems, Inc.
-->
<service_bundle type='manifest' name='opends'>
    <service name='network/ldap/server' type='service' version='1'>
        <!--
            Because we may have multiple instances of network/ldap/server
            provided by different implementations, we keep dependencies
            and methods within the instance.
        -->
        <instance name='opends' enabled='false'>
      	<!--  
     	Wait for network interfaces to be initialized. 
     	-->
            <dependency
       	   name='network' 
       	   grouping='require_all' 
       	   restart_on='error' 
       	   type='service'>
                <service_fmri value='svc:/milestone/network:default' />
            </dependency>
     	<!-- 
       	   Wait for all local filesystems to be mounted. 
     	-->
            <dependency
       	   name='filesystem-local' 
       	   grouping='require_all' 
       	   restart_on='none' 
       	   type='service'>
                <service_fmri value='svc:/system/filesystem/local:default' />
            </dependency>
            <exec_method
            type='method'
            name='start'
            exec='/usr/opends/lib/_svc-opends.sh start'
            timeout_seconds='60'>
                <method_context working_directory='/usr/opends'>
                    <method_credential
                   user='ldap' group='ldap'
privileges='basic,net_privaddr,sys_resource,!proc_info,!file_link_any'
limit_privileges='basic,net_privaddr,sys_resource,!proc_info,!file_link_any'
               />
                </method_context>
            </exec_method>
            <exec_method
            type='method'
            name='stop'
            exec='/usr/opends/lib/_svc-opends.sh stop'
            timeout_seconds='60'>
                <method_context working_directory='/usr/opends'>
                    <method_credential user='ldap' group='ldap' />
                </method_context>
            </exec_method>
            <property_group name='startd' type='framework'>
            <!-- core dumps shouldn't restart session -->
                <propval name='ignore_error' type='astring'
                value='core,signal' />
            </property_group>
            <property_group name='general' type='framework'>
            <!-- autorizations required to start stop opends -->
                <propval name='action_authorization' type='astring'
                    value='solaris.smf.manage.opends' />
                <propval name='value_authorization' type='astring'
                    value='solaris.smf.manage.opends' />
            </property_group>
            <template>
                <common_name>
                    <loctext xml:lang='C'>
                        OpenDS LDAP directory server
                    </loctext>
                </common_name>
                <documentation>
                    <doc_link name='OpenDS Documentation Wiki'
	                          uri='https://www.opends.org/wiki/' />
                </documentation>
            </template>
        </instance>
        <stability value='Evolving' />
    </service>
</service_bundle>
